# config/docker/nginx/files/etc/cron.d/jobs
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

0 12 * * Tue root /certbot.sh >> /var/log/cron.log 2>&1
